# supervisord.conf
[supervisord]
nodaemon=true

[program:migrate]
command=/app/django.sh
autorestart=false

[program:daphne]
command=daphne -b 0.0.0.0 -p 8000 djangotest.asgi:application
autorestart=true
